# srcs/mysql/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mysql]
command= bash mysql.sh
autostart=true
autorestart=true
