# srcs/mysql/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mysql]
command= sh mysql.sh
autostart=true
autorestart=true
